# taptap/rootfs/etc/s6-overlay/s6-rc.d/init-taptap/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Li Tin O've Weedle Add-on: taptap
# Configures TapTap to MQTT bridge
# ==============================================================================

if ! bashio::fs.directory_exists /run/taptap; then
    mkdir -p /run/taptap
    chown -R taptap:taptap /run/taptap
fi

# Generate taptap-mqtt configuration
bashio::var.json \
    mqtt_server "$(bashio::config 'mqtt_server')" \
    mqtt_port "$(bashio::config 'mqtt_port')" \
    mqtt_qos "$(bashio::config 'mqtt_qos')" \
    mqtt_timeout "$(bashio::config 'mqtt_timeout')" \
    mqtt_user "$(bashio::config 'mqtt_user')" \
    mqtt_pass "$(bashio::config 'mqtt_pass')" \
    taptap_serial "$(bashio::config 'taptap_serial')" \
    taptap_address "$(bashio::config 'taptap_address')" \
    taptap_port "$(bashio::config 'taptap_port')" \
    taptap_module_ids "$(bashio::config 'taptap_module_ids')" \
    taptap_module_names "$(bashio::config 'taptap_module_names')" \
    taptap_topic_prefix "$(bashio::config 'taptap_topic_prefix')" \
    taptap_topic_name "$(bashio::config 'taptap_topic_name')" \
    taptap_timeout "$(bashio::config 'taptap_timeout')" \
    taptap_update "$(bashio::config 'taptap_update')" \
    ha_discovery_prefix "$(bashio::config 'ha_discovery_prefix')" \
    ha_birth_topic "$(bashio::config 'ha_birth_topic')" \
    ha_entity_availability "$(bashio::config 'ha_entity_availability')" \
    | tempio \
        -template /etc/taptap/config.gtpl \
        -out /etc/taptap/config.ini
